# roles/podman/templates/mms-image-prune.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Weekly dangling image prune

[Timer]
OnCalendar={{ podman_prune_schedule }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
